# desktop-doorway/etc/systemd/system/desktop-doorway.service
[Unit]
Description=Desktop Doorway
WantedBy=vacant.target

[Service]
ExecStart=/usr/bin/logon-door.sh

[Install]
WantedBy=vacant.target
