# desktop-doorway/etc/systemd/system/user@.service
[Service]
ExecStopPost=+/usr/local/bin/vacancy.sh
