# desktop-doorway/DEBIAN/postinst
sudo systemctl enable vacant.target
sudo systemctl enable desktop-doorway.service
